# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/crond -f -L 15
autorestart=false

[program:server]
command=/http/server
