# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/crond

[program:server]
command=/http/server
